# docker/rootfs/etc/services.d/alice-transmission-dialogs/run
#!/usr/bin/with-contenv sh
alice_gkeep -c /opt/alice-gkeep-dialogs/config.yml
